# systemd/air-collector.service
[Unit]
Description=Pneuma Air Quality Collector
After=network.target influxdb.service

[Service]
Type=simple
User=pneuma
WorkingDirectory=/opt/pneuma
EnvironmentFile=/opt/pneuma/.env
ExecStart=/opt/pneuma/venv/bin/python collector/collector.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
