# systemd/air-api.service
[Unit]
Description=Pneuma Air Quality API
After=network.target influxdb.service

[Service]
Type=simple
User=pneuma
WorkingDirectory=/opt/pneuma
EnvironmentFile=/opt/pneuma/.env
ExecStart=/opt/pneuma/venv/bin/uvicorn api.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
