# systemd/air-lametric.timer
[Unit]
Description=Pneuma LaMetric Push Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
